# backend-hormonia/scripts/webhook-retry.service
[Unit]
Description=Evolution API Webhook Retry Worker
Documentation=file:///app/docs/EVOLUTION_API_WEBHOOK_FIXES.md
After=network.target postgresql.service redis.service
Wants=network-online.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/app/backend-hormonia

# Environment
Environment="PYTHONUNBUFFERED=1"
Environment="WEBHOOK_RETRY_INTERVAL=60"
Environment="ENVIRONMENT=production"
EnvironmentFile=-/etc/hormonia/webhook-retry.env

# Execution
ExecStart=/usr/bin/python3 scripts/webhook_retry_worker.py
ExecReload=/bin/kill -HUP $MAINPID

# Restart policy
Restart=always
RestartSec=10
StartLimitInterval=5min
StartLimitBurst=3

# Resource limits
MemoryLimit=512M
CPUQuota=50%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=webhook-retry

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/app/backend-hormonia/logs

[Install]
WantedBy=multi-user.target
